# nginx-php/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=service nginx start

[program:php5-fpm]
command=/etc/init.d/php5-fpm start
